# src/monea.service
[Unit]
Description=Monea

[Service]
Type=simple
ExecStart=/home/pi/run.sh
User=monea

[Install]
WantedBy=multi-user.target